# codes/x_web/djangoProject/scripts/zipper/test_codes/hitler_adolf/deploy/blackhole.ini
[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[program:thanos]
command=/usr/bin/python3 /opt/stack/dong/codes/hitler_adolf/main.py
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=5
user=root
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/stack/dong/codes/hitler_adolf/logs/cmd_history.log
stderr_logfile=/opt/stack/dong/codes/hitler_adolf/logs/cmd_err.log

[inet_http_server]         ; inet (TCP) server disabled by default
port=0.0.0.0:9001        ; ip_address:port specifier, *:port for all iface
username=admin              ; default is no username (open server)
password=yanfav587
